# services/theia-web.service
[Unit]
Description=THEIA Web (Next.js)
After=network.target theia-api.service
Wants=theia-api.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/theia/app
EnvironmentFile=/opt/theia/app/.env
ExecStart=/usr/bin/node /opt/theia/app/.next/standalone/server.js
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theia-web

[Install]
WantedBy=multi-user.target
